# roles/roon/templates/roonserver.service.j2
# {{ ansible_managed }}
# /etc/systemd/system/roonserver.service
#
# Roon Server systemd unit — managed by Ansible
# Data root: {{ roon_data_root }}

[Unit]
Description=Roon Server
After=network-online.target remote-fs.target
Wants=network-online.target
RequiresMountsFor={{ roon_music_mount }}

[Service]
Type=simple
ExecStart=/opt/RoonServer/start.sh
Environment=ROON_DATAROOT={{ roon_data_root }}
Environment=ROON_ID_DIR=/var/roon
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
LimitMEMLOCK=infinity

# Hardening
ProtectSystem=false
ProtectHome=false
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
